# swiszcli/contrib/systemd/swiszcli-dream.service
[Unit]
Description=swiszCLI dream cycle (nightly chunk promotion + prune + deprecate)
Wants=swiszmem.service
After=swiszmem.service

[Service]
Type=oneshot
ExecStart=%h/.local/bin/swiszcli-dream
Nice=10
IOSchedulingClass=idle
