# swiszcli/contrib/systemd/swiszcli-dream.timer
[Unit]
Description=Run swiszCLI dream cycle nightly at 4am

[Timer]
OnCalendar=*-*-* 04:00:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
